# deployments/fleetd.service
[Unit]
Description=fleetd daemon
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/bin/fleetd
Restart=always
RestartSec=1
User=fleetd
Environment=FLEETD_CONFIG_DIR=/etc/fleetd

[Install]
WantedBy=multi-user.target